# bin/base/files/etc/uci-defaults/999_commit.sh
#!/bin/bash

uci commit network
uci commit firewall
uci commit wireless

# Restart wifi devices
wifi
