# bin/base/files/etc/uci-defaults/999_wifi.sh
#!/bin/bash

# Restart wifi devices
wifi
